# systemd/g4f.service
[Unit]
Description=G4F (GPT4Free) Keyless LLM Aggregator - 12 verified models, 30+ providers
After=network.target

[Service]
Type=simple
ExecStart=__PYTHON_PATH__ -c "from g4f.api import run_api; run_api(port=1337)"
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target